# example/blog-update/systemd/cara.service
[Service]
Type=simple
User=cara
ExecStart=/srv/pipelines/blogupdate.sh
WorkingDirectory=/blogs/cara
