# example/blog-update/systemd/cara.path
[Path]
PathModified=/srv/webhooks/notices/blag_cara

[Install]
WantedBy=default.target
