# example/blog-update/systemd/chuck.path
[Path]
PathModified=/srv/webhooks/notices/blag_chuck

[Install]
WantedBy=default.target
