# example/blog-update/systemd/chuck.service
[Service]
Type=simple
User=chuck
ExecStart=/srv/pipelines/blogupdate.sh
WorkingDirectory=/blogs/chuck
