# example/blog-update/systemd/doug.path
[Path]
PathModified=/srv/webhooks/notices/blag_doug

[Install]
WantedBy=default.target
