# example/blog-update/systemd/doug.service
[Service]
Type=simple
User=doug
ExecStart=/srv/pipelines/blogupdate.sh
WorkingDirectory=/blogs/doug
