# example/blog-update/systemd/fred.path
[Path]
PathModified=/srv/webhooks/notices/blag_fred

[Install]
WantedBy=default.target
