# example/blog-update/systemd/fred.service
[Service]
Type=simple
User=fred
ExecStart=/srv/pipelines/blogupdate.sh
WorkingDirectory=/blogs/fred
