# example/blog-update/systemd/user1.service
[Service]
Type=simple
User=user1
ExecStart=/srv/pipelines/blogupdate.sh
WorkingDirectory=/blogs/user1
