# example/blog-update/systemd/user1.path
[Path]
PathModified=/srv/webhooks/notices/blag_user1

[Install]
WantedBy=default.target
