# systemd/Path.path
[Path]
PathModified=/srv/webhooks/notices/<RepoMap Name>

[Install]
WantedBy=default.target
