# systemd/Path.service
[Service]
Type=oneshot
User=gavin
ExecStart=/srv/daemon/gavin.sh