# systemd/Path.path
[Unit]
Description="Monitor a file for changes, trigger the pipeline."

[Path]
#PathModified doesn't react to touch
PathChanged=/srv/webhooks/notices/TARGET

[Install]
WantedBy=default.target
