# systemd/Path.path
[Unit]
Description="Monitor a file for changes, trigger the pipeline."

[Path]
PathModified=/srv/webhooks/notices/<RepoMap Name>

[Install]
WantedBy=default.target
