# systemd/Path.service
[Service]
Type=simple
User=<username>
ExecStart=/path/to/pipelines/gavin.sh
WorkingDirectory=/path/to/WorkingDirectory