# systemd/Path.path
[Path]
PathModified=/srv/webhooks/waitressd/notices/@user
